# systemd/user/symphony-dashboard-lan-nopilot.socket
[Unit]
Description=LAN listener for Symphony nopilot dashboard

[Socket]
ListenStream=0.0.0.0:4001
NoDelay=true

[Install]
WantedBy=sockets.target
